# puppet/modules/common/files/etc/profile.d/bashrc.sh_dev
# .bashrc
export PS1="\u@\[\e[42m\]\h\[\e[m\] - [$PWD] - \d \A \n> "
