# puppet/modules/mariadb/files/var/spool/cron/crontabs/mysql_prod
* * * * * /usr/bin/mysqldump -uroot --all-databases > /var/lib/mysql/backup.sql 2>/var/lib/mysql/backup.err
